# playbook/roles/tele/files/tele.service
[Unit]
Description=Служба контроля места
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
ExecStart=/root/df.sh

